# supervisord.conf
[supervisord]
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock


[program:webserver]
command = gunicorn --workers=4 --threads=4 --worker-class=gthread --timeout=300 --max-requests=10000 --max-requests-jitter=4 --bind=0.0.0.0:8000 spicebot.wsgi
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup = true


# [program:webserver]
# command = gunicorn --workers 8 --worker-class="egg:meinheld#gunicorn_worker" --timeout 120 --max-requests 1000 --max-requests-jitter 4 --bind 0.0.0.0:8000 spicebot.wsgi
# stdout_logfile=/dev/stdout
# stdout_logfile_maxbytes=0
# stderr_logfile=/dev/stderr
# stderr_logfile_maxbytes=0
# stopasgroup = true


[program:celery_worker_beat]
command = celery -A spicebot worker --beat -l INFO -c 1 -Ofair -Q beat
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup = true


[program:celery_worker_media]
command = celery -A spicebot worker -l INFO -Ofair -Q media --autoscale=1,3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup = true


[program:celery_worker_deposit]
command = celery -A spicebot worker -l INFO -c 1 -Ofair -Q deposit --max-tasks-per-child=1 
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup = true


[program:celery_worker_negative_balance]
command = celery -A spicebot worker -l INFO -c 1 -Ofair -Q negative_balance --max-tasks-per-child=1
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup = true

# [program:celery_worker_deposit_socket]
# command = celery -A spicebot worker -l INFO -c 3 -Ofair -Q deposit_socket
# stdout_logfile=/dev/stdout
# stdout_logfile_maxbytes=0
# stderr_logfile=/dev/stderr
# stderr_logfile_maxbytes=0
# stopasgroup = true


# [program:celery_worker_deposit_block]
# command = celery -A spicebot worker -l INFO -c 2 -Ofair -Q deposit_block
# stdout_logfile=/dev/stdout
# stdout_logfile_maxbytes=0
# stderr_logfile=/dev/stderr
# stderr_logfile_maxbytes=0
# stopasgroup = true


[program:celery_worker_withdrawal]
command = celery -A spicebot worker -l INFO -c 1 -Ofair -Q withdrawal --max-tasks-per-child=1
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup = true


[program:celery_worker_telegram]
command = celery -A spicebot worker -l INFO -Ofair -Q telegram --autoscale=1,8 --max-tasks-per-child=100
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup = true


# [program:celery_worker_twitter]
# command = celery -A spicebot worker -l INFO -c 2 -Ofair -Q twitter
# stdout_logfile=/dev/stdout
# stdout_logfile_maxbytes=0
# stderr_logfile=/dev/stderr
# stderr_logfile_maxbytes=0
# stopasgroup = true

# [program:celery_worker_twitter_dm]
# command = celery -A spicebot worker -l INFO -c 1 -Ofair -Q twitter_dm
# stdout_logfile=/dev/stdout
# stdout_logfile_maxbytes=0
# stderr_logfile=/dev/stderr
# stderr_logfile_maxbytes=0
# stopasgroup = true

# [program:celery_worker_reddit]
# command = celery -A spicebot worker -l INFO -c 1 -Ofair -Q reddit
# stdout_logfile=/dev/stdout
# stdout_logfile_maxbytes=0
# stderr_logfile=/dev/stderr
# stderr_logfile_maxbytes=0
# stopasgroup = true

[program:celery_worker_misc]
command = celery -A spicebot worker -l INFO -c 2 -Ofair -Q faucet,transfer,bridge --max-tasks-per-child=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup = true

# [program:celery_worker_transfer]
# command = celery -A spicebot worker -l INFO -c 1 -Ofair -Q transfer --max-tasks-per-child=10
# stdout_logfile=/dev/stdout
# stdout_logfile_maxbytes=0
# stderr_logfile=/dev/stderr
# stderr_logfile_maxbytes=0
# stopasgroup = true

# [program:subreddit]
# command = celery -A spicebot worker -l INFO -c 1 -Ofair -Q subreddit
# stdout_logfile=/dev/stdout
# stdout_logfile_maxbytes=0
# stderr_logfile=/dev/stderr
# stderr_logfile_maxbytes=0
# stopasgroup = true

[program:celery_error_logs]
command = celery -A spicebot worker -l INFO -c 1 -Ofair -Q logs --max-tasks-per-child=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup = true

# Temporarily disabled.
# [program:celery_spice_web_hook]
# command = celery -A spicebot worker -l INFO -c 4 -Ofair -Q spice_web_hook
# stdout_logfile=/dev/stdout
# stdout_logfile_maxbytes=0
# stderr_logfile=/dev/stderr
# stderr_logfile_maxbytes=0
# stopasgroup = true

# [program:celery_missed_deposit_block]
# command = celery -A spicebot worker -l INFO -c 1 -Ofair -Q missed_deposit_block
# stdout_logfile=/dev/stdout
# stdout_logfile_maxbytes=0
# stderr_logfile=/dev/stderr
# stderr_logfile_maxbytes=0
# stopasgroup = true

# [program:celery_per_block_deposit_scanner]
# command = celery -A spicebot worker -l INFO -c 8 -Ofair -Q per_block_deposit_scanner
# stdout_logfile=/dev/stdout
# stdout_logfile_maxbytes=0
# stderr_logfile=/dev/stderr
# stderr_logfile_maxbytes=0
# stopasgroup = true

# [program:celery_scan_user_deposits]
# command = celery -A spicebot worker -l INFO -c 8 -Ofair -Q scan_user_deposits
# stdout_logfile=/dev/stdout
# stdout_logfile_maxbytes=0
# stderr_logfile=/dev/stderr
# stderr_logfile_maxbytes=0
# stopasgroup = true

# [program:celery_worker_mute_expiration]
# command = celery -A spicebot worker -l INFO -c 2 -Ofair -Q mute_expiration
# stdout_logfile=/dev/stdout
# stdout_logfile_maxbytes=0
# stderr_logfile=/dev/stderr
# stderr_logfile_maxbytes=0
# stopasgroup = true

[program:celery_worker_record_metrics]
command = celery -A spicebot worker -l INFO -c 2 -Ofair -Q metrics --max-tasks-per-child=1
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopasgroup = true

# [program:celery_worker_process_wagerbot_bet]
# command = celery -A spicebot worker -l INFO -c 2 -Ofair -Q process_wagerbot_bet
# stdout_logfile=/dev/stdout
# stdout_logfile_maxbytes=0
# stderr_logfile=/dev/stderr
# stderr_logfile_maxbytes=0
# stopasgroup = true
